# itsm-system.service
[Unit]
Description=ITSM-System - IT Service Management & NIST CSF 2.0 Compliance System
Documentation=https://github.com/Kensan196948G/ITSM-System
After=network.target
Wants=itsm-frontend.service

[Service]
Type=simple
User=kensan
WorkingDirectory=/mnt/LinuxHDD/ITSM-System
EnvironmentFile=/mnt/LinuxHDD/ITSM-System/.env

# ExecStart - プロダクション環境でのNode.js直接起動
ExecStart=/usr/bin/node backend/server.js

# 自動再起動設定
Restart=always
RestartSec=10

# プロセス管理
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# リソース制限
LimitNOFILE=65536
LimitNPROC=4096

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=itsm-system

[Install]
WantedBy=multi-user.target
