# services/itsm-sec-nexus-dev.service
[Unit]
Description=ITSM-Sec Nexus - 開発環境 (HTTP)
Documentation=https://github.com/yourusername/itsm-sec-nexus
After=network.target

[Service]
Type=simple
User=kensan
Group=kensan
WorkingDirectory=/mnt/LinuxHDD/ITSM-System

# 開発環境用の環境変数ファイル
EnvironmentFile=/mnt/LinuxHDD/ITSM-System/config/env/.env.development

# Start command
ExecStart=/usr/bin/node /mnt/LinuxHDD/ITSM-System/backend/server.js

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=append:/mnt/LinuxHDD/ITSM-System/logs/backend-dev.log
StandardError=append:/mnt/LinuxHDD/ITSM-System/logs/backend-dev.log
SyslogIdentifier=itsm-sec-nexus-dev

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
