# systemd/itsm-system-https.service
[Unit]
Description=ITSM-System Backend (HTTPS)
Documentation=https://github.com/Kensan196948G/ITSM-System
After=network.target network-online.target
Wants=network-online.target itsm-frontend-https.service

[Service]
Type=simple
User=kensan
Group=kensan
WorkingDirectory=/mnt/LinuxHDD/ITSM-System

# 環境変数ファイル
EnvironmentFile=/mnt/LinuxHDD/ITSM-System/.env

# HTTPSサーバー起動
ExecStart=/usr/bin/node backend/server.js

# 起動前のヘルスチェック（SSL証明書の確認）
ExecStartPre=/usr/bin/test -f /mnt/LinuxHDD/ITSM-System/ssl/server.crt
ExecStartPre=/usr/bin/test -f /mnt/LinuxHDD/ITSM-System/ssl/server.key

# 自動再起動設定
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# プロセス管理
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=60
TimeoutStopSec=30

# リソース制限
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=2G
CPUQuota=200%

# セキュリティ設定（強化版）
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true

# 読み書き可能なパス
ReadWritePaths=/mnt/LinuxHDD/ITSM-System/backend
ReadWritePaths=/mnt/LinuxHDD/ITSM-System/ssl

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=itsm-system-https

# 環境変数の上書き（HTTPSを強制的に有効化）
Environment="ENABLE_HTTPS=true"
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target
