# systemd/itsm-system.service
[Unit]
Description=ITSM-Sec Nexus System (Native Node.js)
Documentation=https://github.com/your-org/itsm-system
After=network.target

[Service]
Type=simple
User=kensan
Group=kensan
WorkingDirectory=/mnt/LinuxHDD/ITSM-System

# Environment
Environment="NODE_ENV=production"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=/mnt/LinuxHDD/ITSM-System/.env.production

# Start command
ExecStart=/usr/bin/node backend/server.js

# Restart policy
Restart=always
RestartSec=10s

# Resource limits
MemoryLimit=2G
CPUQuota=200%

# Logging
StandardOutput=append:/mnt/LinuxHDD/ITSM-System/logs/itsm-system.log
StandardError=append:/mnt/LinuxHDD/ITSM-System/logs/itsm-system-error.log
SyslogIdentifier=itsm-system

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
